# masakari-engine.service
[Unit]
Description=Masakari Engine Daemon

[Service]
Type=simple
User=openstack
Group=openstack
ExecStart=/usr/bin/masakari-engine

[Install]
WantedBy=multi-user.target
