# masakari-api.service
[Unit]
Description=Masakari Api Daemon

[Service]
Type=simple
User=masakari
ExecStart=/usr/bin/masakari-api

[Install]
WantedBy=multi-user.target
